# src/supervisord.conf
[supervisord]
nodaemon=true

[program:fpm]
command=php-fpm -F

[program:nginx]
command=nginx

